# roles/data/templates/backup-grafana-db.service
[Unit]
Description=Backup postgres data

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-grafana-db
