# roles/data/templates/backup-postgres-data.service
[Unit]
Description=prune-influx-data

[Service]
User=postgres
Type=oneshot
ExecStart=/usr/local/bin/backup-postgres-data
