# roles/data/templates/prune-influx-data.service
[Unit]
Description=prune-influx-data

[Service]
Type=oneshot
ExecStart=/usr/local/bin/prune-influx-data
